# marcelo-scripts/kiev/fwkiev_team-git.cron
#every 7min
*/7 * * * * gitdummy /usr/local/bin/fwkiev_team_git
