# marcelo-scripts/kiev/push_denver.cron
#every 7min
*/7 * * * * gitdummy /usr/local/bin/push_denver
